# docker/etc/s6-overlay/s6-rc.d/init-fixuser/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

groupmod -o -g "${PGID}" ab
usermod -o -u "${PUID}" ab

chown ab:ab -R /app /ab
