# docker/etc/s6-overlay/s6-rc.d/init-old-compatible/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask ${UMASK}

if [ -f /config/bangumi.json ]; then
    mv /config/bangumi.json /app/data/bangumi.json
fi
