# rootfs/etc/services.d/ab_api/run
#!/usr/bin/with-contenv bash

cd /src
umask ${UMASK}
s6-setuidgid ${PUID}:${PGID} python3 api.py